# mpv_launcher.service
[Unit]
Description=mpv_launcher

[Service]
ExecStart=/bin/bash /home/teacher/Документы/bash_mpv_launcher.sh

[Install]
WantedBy=default.target

